# dockerfiles/code/root/etc/services.d/code/run
#!/usr/bin/with-contenv bash

set -e

s6-setuidgid $USER code $CODE_ARGS

while true; do 
    pgrep -f /usr/share/code/code &>/dev/null
    sleep 10
done