# dockerfiles/ubuntu-x11/root/etc/services.d/xfdesktop/run
#!/usr/bin/with-contenv bash

set -e
if [[ -n $XFDESKTOP ]]; then
    s6-setuidgid $USER xfdesktop $XFDESKTOP_PANEL_ARGS
else
    sleep infinity
fi

sleep 5