# dockerfiles/ubuntu/root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash


crontab /etc/cron.d/*

/usr/sbin/cron -f -l -L 15