# dockerfiles/vivaldi-novnc/root/etc/services.d/vivaldi/run
#!/usr/bin/with-contenv bash

set -e

s6-setuidgid $USER vivaldi $VIVALDI_ARGS

while true; do 
    pgrep vivaldi &>/dev/null
    sleep 10
done