# dockerfiles/pycharm/root/etc/services.d/pycharm/run
#!/usr/bin/with-contenv bash

set -e

s6-setuidgid $USER /opt/pycharm/bin/pycharm.sh

while true; do 
    pgrep -f pycharm.sh &>/dev/null
    sleep 10
done