# dockerfiles/code-server/root/etc/services.d/rsync/run
#!/usr/bin/with-contenv bash

while true; do
    s6-setuidgid $USER rsync -hrtp /home/coder/.local/share/code-server/ /home/coder/.local/share/code-server-host/
    sleep 3
done