# dockerfiles/novnc/root/etc/services.d/xterm/run
#!/usr/bin/with-contenv bash

xhost +
s6-setuidgid $USER xterm