# dockerfiles/phpstorm/root/etc/services.d/phpstorm/run
#!/usr/bin/with-contenv bash

set -e

s6-setuidgid $USER /opt/phpstorm/bin/phpstorm.sh

while true; do 
    pgrep -f phpstorm.sh &>/dev/null
    sleep 10
done