# dockerfiles/novnc-docker/root/etc/services.d/dind/run
#!/usr/bin/with-contenv bash

rm -f /var/run/docker.pid
dind dockerd